# supervisor.conf
[supervisord]
nodaemon=true

[program:chat]
directory=/home/app
command=/home/app/chat
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/log/supervisor/chat.log
